# Mydomohome-plugins/docker-owhttp/supervisor.owhttpd.conf
[program:owhttpd]
command=owhttpd -u -p 80
directory=/
autostart=true
autorestart=true
logfile=/var/log/supervisor/owhttpd.log
log_stderr=true
user=root

